# config/launchagents/com.urfael.morningbrief.plist.template
<?xml version="1.0" encoding="UTF-8"?>
<!-- Urfael morning brief — opt-in. install.sh fills {{HOME}}. Change Hour to retime. -->
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0"><dict>
  <key>Label</key><string>com.urfael.morningbrief</string>
  <key>ProgramArguments</key>
  <array><string>/bin/bash</string><string>{{HOME}}/Urfael/_urfael/morning-brief.sh</string></array>
  <key>StartCalendarInterval</key><dict><key>Hour</key><integer>8</integer><key>Minute</key><integer>0</integer></dict>
  <key>RunAtLoad</key><false/>
  <key>EnvironmentVariables</key>
  <dict><key>PATH</key><string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string></dict>
  <key>StandardOutPath</key><string>{{HOME}}/.claude/urfael/morningbrief.out.log</string>
  <key>StandardErrorPath</key><string>{{HOME}}/.claude/urfael/morningbrief.err.log</string>
</dict></plist>
